# scripts/pimonervix.service
[Unit]
Description=pimonervix

[Service]
ExecStart=/var/www/pimonervix/scripts/start
User=pi

[Install]
WantedBy=multi-user.target
